# systemd/pico-cec-boot.service
# /etc/systemd/system/pico-cec-boot.service
#
# Fires once on system boot (cold boot, not resume -- resume is handled by
# the systemd-sleep hook) to wake the TV and switch it to this PC's input.
#
# Bound to graphical.target: on a HTPC this fires when the display manager
# / desktop is ready. (graphical-session.target is a user-session target not
# available to system units.)

[Unit]
Description=Pico CEC: TV power-on at boot
After=pico-cec-listener.service graphical.target
Requires=pico-cec-listener.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pico-cec-ctl.py PWR_ON
# Retry a couple of times with backoff: right at boot the Pico may still
# be finishing its own EDID read / CEC bus settle, so an immediate PWR_ON
# could race ahead of EVT:READY. pico-cec-ctl already blocks up to ~8s for
# an ACK per call; this wraps that in a couple of outer attempts too.
ExecStartPre=/bin/sleep 2

[Install]
WantedBy=graphical.target
