# systemd/pico-cec-boot.service
# /etc/systemd/system/pico-cec-boot.service
#
# Fires once on system boot (cold boot, not resume -- resume is handled by
# the systemd-sleep hook) to wake the TV and switch it to this PC's input.
#
# Bound to graphical.target: on a HTPC this fires when the display manager
# / desktop is ready. (graphical-session.target is a user-session target not
# available to system units.)

[Unit]
Description=Pico CEC: TV power-on at boot
After=pico-cec-listener.service graphical.target
Requires=pico-cec-listener.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pico-cec-ctl.py PWR_ON
# Small startup delay: give the listener and USB serial path a moment to
# settle before the initial PWR_ON.
ExecStartPre=/bin/sleep 2

[Install]
WantedBy=graphical.target
