# systemd/pico-cec-listener.service
# /etc/systemd/system/pico-cec-listener.service
#
# Long-running daemon: owns the serial connection to the Pico, reacts to
# TV-standby events by suspending the system, and exposes a control socket
# at /run/pico-cec/control.sock for the sleep/shutdown hooks to request
# PWR_ON / PWR_OFF without fighting over the serial port.

[Unit]
Description=Pico CEC bridge listener
After=multi-user.target
Wants=dev-pico\x2dcec.device
After=dev-pico\x2dcec.device

[Service]
Type=simple
ExecStart=/usr/local/bin/pico-cec-listener.py
Restart=always
RestartSec=2
RuntimeDirectory=pico-cec
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
